# podman/supervisord.conf
[supervisord]
nodaemon=true

[program:vncserver]
command=/usr/bin/vncserver
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

